# packaging/api/root/usr/share/utisak-api/init/utisak-api.conf
start on (filesystem and net-device-up)

chdir /var/lib/utisak-api
console log

script
  /usr/local/bin/utisak-api --config=/etc/utisak-api/config.json > /var/log/utisak-api.log  2>&1
end script
