# packaging/worker/root/usr/share/utisak-worker/init.d/utisak-worker
#! /bin/sh

### BEGIN INIT INFO
# Provides:          utisak-worker
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: utisak-worker
### END INIT INFO

DAEMON_OPTS="--config=/etc/utisak-worker/config.json"

pid() {
  if [ -f /usr/local/bin/utisak-worker ]; then
    pidof /usr/local/bin/utisak-worker
  fi
}

stop() {
  if pidof /usr/local/bin/utisak-worker >/dev/null; then
    kill -9 "$(pid)"
  else
    echo "utisak-worker not runned"
    exit 1
  fi
}

start() {
  if pidof /usr/local/bin/utisak-worker >/dev/null; then
    echo "utisak-worker already runned"
    exit 1
  else
    nohup utisak-worker $DAEMON_OPTS > /var/log/utisak-worker.log  2>&1 &
  fi
}

restart() {
  if pidof /usr/local/bin/utisak-worker >/dev/null; then
    kill -9 "$(pid)"
    nohup utisak-worker $DAEMON_OPTS > /var/log/utisak-worker.log  2>&1 &
    exit 0
  else
    nohup utisak-worker $DAEMON_OPTS > /var/log/utisak-worker.log  2>&1 &
    exit 0
  fi
}

status() {
  if pidof /usr/local/bin/utisak-worker >/dev/null; then
    echo "utisak-worker with pid $(pid) is running"
  else
    echo "utisak-worker is not running"
  fi
  exit 0
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  status)
        status
        ;;
  *)
        echo "Usage: service utisak-worker {start|stop|restart|status}"
        exit 1
        ;;
esac
